# photobox.service
[Unit]
Description=Photobox Flask Application
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/photobox
Environment=PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10

# Umgebungsvariablen
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1

# Sicherheit
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/photobox

[Install]
WantedBy=multi-user.target